# ifconfig.socket
[Socket]
ListenStream = 8000
BindIPv6Only = both

[Install]
WantedBy = sockets.target
